# deploy/ugc.service
[Unit]
Description=UGC Machine
After=network.target

[Service]
WorkingDirectory=/home/user/ugc-machine
ExecStart=/usr/bin/python3 /home/user/ugc-machine/server.py
Restart=always
User=user

[Install]
WantedBy=multi-user.target
